# systemd/publish-gh-pages.timer
[Unit]
Description=Publish GitHub Pages every 10 minutes

[Timer]
OnCalendar=*:0/10
RandomizedDelaySec=60
Persistent=true

[Install]
WantedBy=timers.target
