# systemd/speedtest.timer
[Unit]
Description=Speedtest every 10 minutes

[Timer]
OnCalendar=*:0/10
RandomizedDelaySec=30
Persistent=true

[Install]
WantedBy=timers.target
